# spyder-setup/mounthdd.service
[Unit]
Description=Mount HDD

[Service]
Type=oneshot
ExecStart=/home/pepper/mounthdd.sh
