# spyder-setup/mounthdd.timer
[Unit]
Description=Run mounthdd.service at boot

[Timer]
OnBootSec=1sec

[Install]
WantedBy=timers.target
